# modules/install-k8s/k8s.target
[Unit]
Description=K8s target
After=cloud-init.target
